# btw_i_use_arch/config/rest_uploader.ini
[program:rest_uploader]
command=/home/justin/.virtualenvs/rest_uploader/bin/python -m rest_uploader.cli -d inbox /home/justin/.joplin_upload --moveto /home/justin/.joplin_upload/imported
directory=/home/justin
autostart=true
autorestart=true
stderr_logfile=/var/log/rest_uploader.err.log
stdout_logfile=/var/log/rest_uploader.out.log
user=justin
group=justin
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,PYTHONUNBUFFERED=1
